# repo/pkg-g/g-ad/pkg.install
#!/bin/bash

post_install() {
	systemctl enable systemd-timesyncd.service
}