# repo/pkg-g/g-ad/pkg.install
#!/bin/bash

post_install() {
    systemctl disable systemd-resolved
    systemctl enable smb
    systemctl enable nmb
}