# repo/pkg-g/g-core/pkg.install
#!/bin/bash

post_install() {
    systemctl enable NetworkManager
    systemctl enable sshd
    systemctl enable goylin.service
}