# repo/pkg-g/g-backgrounds/pkg.install
#!/bin/bash

post_install() { systemctl enable goylinGenWall.service; }

pre_upgrade() { [ -d /usr/share/goylin/backgrounds ] && rm /usr/share/goylin/backgrounds/*; }