# repo/pkg-g/g-printer/pkg.install
#!/bin/bash

pre_upgrade() {
	chattr -i /etc/cups/printers.conf
}

post_install() {
	systemctl enable cups
	chattr +i /etc/cups/printers.conf
}

post_upgrade() {
	chattr +i /etc/cups/printers.conf
}