# repo/pkg-g/g-printer/pkg.install
#!/bin/bash

post_install() {
	systemctl enable cups
	systemctl enable cups-browsed
}