# repo/pkg-g/g-desktop/skel/.config/autostart/goylinWelcome.desktop
[Desktop Entry]
Type=Application
Exec=/usr/bin/gWelcome
Hidden=false
NoDisplay=false
X-GNOME-Autostart-enable=true
Name=GoylinWelcome