# repo/pkg-g/g-desktop/skel/.config/autostart/nm-applet.desktop
[Desktop Entry]
Type=Application
Icon=application-x-shellscript
Name=nm-applet
Exec=/usr/bin/nm-applet
X-KDE-AutostartScript=true
