# repo/pkg-g/g-desktop/skel/.config/autostart/screenResolution.desktop
[Desktop Entry]
Type=Application
Icon=application-x-shellscript
Name=gSetResolution
Exec=/usr/bin/gSetResolution
X-KDE-AutostartScript=true
