# repo/pkg-g/g-desktop/skel/.config/autostart/screenResolution.desktop
[Desktop Entry]
Type=Application
Exec=/usr/bin/gSetResolution
Hidden=true
NoDisplay=true
X-GNOME-Autostart-enable=true
Name=gSetResolution