# repo/pkg-g/g-desktop/gShutdown.service
[Unit]
Description=Shutdown system on Time

[Service]
Type=oneshot
ExecStart=/usr/bin/shutdown -h now