# repo/pkg-g/g-desktop/pkg.install
#!/bin/bash

post_install() {
	systemctl enable lightdm
	chattr +i /var/lib/lightdm/.cache/slick-greeter/state
}

pre_upgrade() {
	chattr -i /var/lib/lightdm/.cache/slick-greeter/state
}

post_upgrade() {
	chattr +i /var/lib/lightdm/.cache/slick-greeter/state
}