# repo/pkg-g/g-desktop/pkg.install
#!/bin/bash

_lockGFiles() { chattr +i /etc/profile; }

_unlockGFiles() { chattr -i /etc/profile; }

post_install() {
	printf "\nsource /etc/profile.goylin\n" | tee -a /etc/profile 1> /dev/null
	
	_lockGFiles

	systemctl enable bluetooth
}

pre_upgrade() { _unlockGFiles; }

post_upgrade() { _lockGFiles; }