# repo/pkg-g/g-desktop/pkg.install
#!/bin/bash

post_install() {
    systemctl enable lightdm
    sudo chattr +i /var/lib/lightdm/.cache/slick-greeter/state
}